# backend/machine_images/roles/datadog-agent/files/etc/systemd/system/datadog-agent.service.d/override.conf
# Note that these overrides are merged with default After/Wants for datadog-agent.service
[Unit]
After=cloud-final.service
Wants=cloud-final.service
